# overlay/etc/profile.d/costapine.sh
#!/bin/sh
# Costapine Linux - Shell Profile
export EDITOR=mousepad
export VISUAL=mousepad
export PAGER=less
export BROWSER=badwolf

alias ll='ls -lah --color=auto'
alias la='ls -A --color=auto'
alias l='ls -CF --color=auto'
alias grep='grep --color=auto'
